# deploy/digital-nameplate.service
[Unit]
Description=Digital Office Nameplate System
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/digital-nameplate
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=5
Environment=NODE_ENV=production
Environment=PORT=3004
EnvironmentFile=-/opt/digital-nameplate/.env

[Install]
WantedBy=multi-user.target
